# api.service
[Unit]
Description=API

[Service]
ExecStart= /usr/bin/node /home/terraform/TeamAPoolProjectBackend/Index.js

[Install]
WantedBy=multi-user.target
