# api.service
[Unit]
Description=API

[Service]
ExecStart= node ~/TeamAPoolProjectBackend/Index.js

[Install]
WantedBy=multi-user.target
