# api.service
[Unit]
Description=API

[Service]
ExecStart= node /home/allancharles2095/TeamAPoolProjectBackend/Index.js

[Install]
WantedBy=multi-user.target
