# ui.service
[Unit]
Description=ui

[Service]
WorkingDirectory= ~/TeamAPoolProjectUI
ExecStart= ng s --host 0.0.0.0

[Install]
WantedBy=mult-user.target
